# .config/autostart/tint2.desktop
[Desktop Entry]
Name=Tint2 taskbar
Comment=Tint2 taskbar
Exec=tint2
Terminal=false
Type=Application
StartupNotify=true
X-GNOME-Autostart-enabled=true